# vscode-ssh/99-usb-permissions.rules
# Allow dev user to access all USB devices
SUBSYSTEM=="usb", MODE="0666"
SUBSYSTEM=="usb_device", MODE="0666"

# Allow mounting of external media
SUBSYSTEM=="block", ATTR{removable}=="1", ENV{ID_FS_USAGE}=="filesystem", ACTION=="add", RUN+="/bin/mount -t auto -o rw /dev/%k /media/%k"

# Grant access to disk devices
KERNEL=="sd[a-z]*", MODE="0666"

# Allow access to serial devices
KERNEL=="ttyUSB[0-9]*", MODE="0666"
KERNEL=="ttyACM[0-9]*", MODE="0666"
